# jcm.service
[Unit]
Description=JCM Clipboard Manager Daemon
After=graphical-session.target

[Service]
ExecStart=jcm-daemon watch
Restart=on-failure

[Install]
WantedBy=graphical-session.target
